# bot.conf
[program:KPI-Rozklad-Bot]
directory = /home/ubuntu/KPI-Rozklad-Bot
command = python3 app.py



autostart = true
autorestart = true

enviroment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile = home/ubuntu/logfile_err.log
stdout_logfile = home/ubuntu/logfile.log